# deb-src/uucin-tomcat7-bottle-libGeomlibrary/etc/profile.d/uucin-tomcat7-bottle-libGeomlibrary.sh
LD_LIBRARY_PATH=/usr/local/uu_tomcat7_bottle/lib:$LD_LIBRARY_PATH
export LD_LIBRARY_PATH
